# entries/kql/0199bf37-eb5f-73c4-a514-5f5daeb561ab.yaml
id: 0199bf37-eb5f-73c4-a514-5f5daeb561ab
name: Scheduled Task Leading to DLL Sideloading from Suspicious Paths
description: Detects the creation of a scheduled task that launches an executable from a non-standard, writable directory, which subsequently loads a co-located DLL. This pattern is a form of DLL Sideloading used for persistence and was observed in a ToddyCat APT attack.
query: "// Name: Scheduled Task leading to DLL Sideloading\r\n// Author: Zelfatek powered by detections.ai\r\n// Date: 2025-10-07\r\n// Tactic: Persistence, Defense Evasion, Privilege Escalation\r\n// Technique: T1053.005 (Scheduled Task), T1574.001 (DLL Sideloading)\r\n// Description: Detects the creation of a scheduled task that launches an executable from a non-standard, writable directory, which subsequently loads a co-located DLL. This pattern is a form of DLL Sideloading used for persistence and was observed in a ToddyCat APT attack.\r\n// References: https://securelist.com/detecting-dll-hijacking-with-machine-learning-in-kaspersky-siem\r\n// False Positive Sensitivity: Medium. Legitimate software installers or updaters might use this pattern. Tune by adding known legitimate executables to an exclusion list.\r\n\r\nlet lookback = 15m;\r\n// Define suspicious, commonly world-writable paths where applications shouldn't typically be installed for persistence.\r\nlet suspiciousPaths = dynamic([\r\n    @\"C:\\ProgramData\\\",\r\n    @\"C:\\Users\\Public\\\",\r\n    @\"C:\\PerfLogs\\\",\r\n    @\"C:\\Windows\\Temp\\\",\r\n    @\"C:\\Temp\\\"\r\n]);\r\n// Find processes launched from suspicious locations that load a DLL from the same directory (sideloading).\r\nlet SideloadEvents =\r\n    DeviceImageLoadEvents\r\n    | where Timestamp > ago(lookback)\r\n    | where FileName endswith \".dll\" and isnotempty(InitiatingProcessFolderPath)\r\n    // The core sideloading condition: process path is the same as the loaded module path.\r\n    | where InitiatingProcessFolderPath == FolderPath\r\n    // Focus on suspicious, non-standard paths.\r\n    | where InitiatingProcessFolderPath has_any (suspiciousPaths)\r\n    | project\r\n        SideloadTime = Timestamp,\r\n        DeviceId,\r\n        DeviceName,\r\n        SideloadingProcess = InitiatingProcessFileName,\r\n        SideloadingProcessPath = InitiatingProcessFolderPath,\r\n        LoadedDll = FileName,\r\n        ReportId,\r\n        InitiatingProcessId;\r\n// Find scheduled task creations that point to executables in these suspicious locations.\r\nlet SchTaskCreations =\r\n    DeviceProcessEvents\r\n    | where Timestamp > ago(lookback * 2) // Look back a bit further for the task creation.\r\n    | where FileName =~ \"schtasks.exe\"\r\n    | where ProcessCommandLine has \"/create\" and ProcessCommandLine has \"/tr\"\r\n    // Extract the executable path from the /tr argument. Handles both quoted and unquoted paths.\r\n    | extend QuotedTaskPath = extract(\"(?i)/tr\\\\s+\\\"([^\\\"]+)\\\"\", 1, ProcessCommandLine)\r\n    | extend UnquotedTaskPath = extract(\"(?i)/tr\\\\s+([^\\\\s\\\"]+)\", 1, ProcessCommandLine)\r\n    | extend TaskExecutablePath = iif(isnotempty(QuotedTaskPath), QuotedTaskPath, UnquotedTaskPath)\r\n    | where isnotempty(TaskExecutablePath)\r\n    // Focus on tasks pointing to executables in the same suspicious paths.\r\n    | where TaskExecutablePath has_any (suspiciousPaths)\r\n    | project\r\n        SchTaskTime = Timestamp,\r\n        DeviceId,\r\n        TaskCreator = InitiatingProcessFileName,\r\n        SchtaskCommandLine = ProcessCommandLine,\r\n        TaskExecutablePath;\r\n// Join the two event types to find scheduled tasks that lead to DLL sideloading.\r\nSideloadEvents\r\n| join kind=inner (\r\n    SchTaskCreations\r\n) on DeviceId\r\n// Match the executable path from the scheduled task with the full path of the process that performed the sideload.\r\n| where strcat(SideloadingProcessPath, SideloadingProcess) == TaskExecutablePath\r\n// Ensure the task was created before the sideloading event occurred, within our lookback window.\r\n| where SideloadTime > SchTaskTime and (SideloadTime - SchTaskTime) < lookback\r\n| project\r\n    Timestamp = SideloadTime,\r\n    DeviceName,\r\n    SideloadingProcess,\r\n    SideloadingProcessPath,\r\n    LoadedDll,\r\n    SchtaskCommandLine,\r\n    TaskCreator,\r\n    SchTaskTime\r\n"
query_type: kql
severity: medium
tactics:
    - TA0002
    - TA0003
    - TA0004
    - TA0005
techniques:
    - T1053.005
    - T1574.001
tags:
    - kql
version: "2"
references:
    - https://securelist.com/detecting-dll-hijacking-with-machine-learning-in-kaspersky-siem
