# entries/kql/019a019d-ec59-7290-9442-3b58368ce798.yaml
id: 019a019d-ec59-7290-9442-3b58368ce798
name: Xubuntu.org Watering Hole Attack Crypto Drainer Malware Detection
description: Detects indicators associated with a watering hole attack originating from a compromised xubuntu.org website. This includes the creation of a tampered 'tos.txt' file, the dropping and execution of a trojanized dropper ('elzvcf.exe') in AppData\Roaming, and the establishment of persistence via a registry Run key.
query: "// Name: Xubuntu.org Website Compromised via Watering Hole Attack Crypto Drainer Malware\n// Author: Reyben T. Cortes @ Pipeline Horizon, Unit Zero\n// Date: 2025-10-20\n// Description: Detects indicators associated with a watering hole attack originating from a compromised xubuntu.org website. This includes the creation of a tampered 'tos.txt' file, the dropping and execution of a trojanized dropper ('elzvcf.exe') in AppData\\Roaming, and the establishment of persistence via a registry Run key.\n// Reference: https://cybernews.com/security/xubuntu-site-compromise-hackers-peddle-malware\n// Tactic: Execution, Persistence, Initial Access\n// Technique: T1204.002 (Malicious File Execution), T1547.001 (Registry Run Keys), T1189 (Drive-by Compromise)\n// False Positive Sensitivity: Medium\n\n// Malicious file hashes associated with Xubuntu Compromise & Crypto Drainer Malware\nlet MaliciousHashes = dynamic([\n    \"38eaef9fee4585f112dd35d5108f9f3baa1977cb682a7068a90996f166c840f5\", // qjrqar.exe\n    \"afaebc6cf20f32ea0644f69c511a5da12f3b860f7d13b18500051830337965d7\", // Stage 2 (elzvcf.exe)\n    \"ec3a45882d8734fcff4a0b8654d702c6de8834b6532b821c083c1591a0217826\", // Stage 1 (TestCompany.SafeDownloader.exe)\n    \"0f59f553fcfac3cac07aa7986eac914be069a6dd407b2d9f761f11d3e865b4f6\", // Malicious ZIP\n    \"4eeb4a1d56b38be699c9c379de5d7d4291f656b766ae21c43843af3ddf55d7e9\", // CorporateDownloadTool.exe\n    \"b09ec204f734363b8826027da5c818b0925ade50d1a645a44d04bb1c4d209a83\"  // Xubuntu_Multi_Download_Manager.zip\n]);\n\n// Look for malicious dropped files starting with .iso torrent and dropped .zip file containing .exe and tampered tos.txt\nlet FileEvents = DeviceFileEvents\n| where Timestamp > ago(7d)\n| where (SHA256 == \"bfaf3af21644ce06fb295c78e080477d9095cb185fcf164ce9cf066a9682b9f1\" and FileName == \"tos.txt\" and FolderPath has_any (\"Downloads\")) \n    or (SHA256 == \"afaebc6cf20f32ea0644f69c511a5da12f3b860f7d13b18500051830337965d7\" or FileName =~ \"elzvcf.exe\" and FolderPath has \"AppData\" and FolderPath has \"Roaming\")  // Payload drop in %APPDATA%\\Roaming\n    or (SHA256 in (MaliciousHashes) and FileOriginRefererUrl has \"xubuntu.org\")  // Download origin from compromised site\n| extend ActivityType = \"File\", EntityType = \"File\", ActionType = \"FileCreateOrDownload\", TargetPath = FolderPath;\n\n// Look for the execution of the known malicious payload from the AppData\\Roaming directory\nlet ProcessEvents = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where SHA256 in (MaliciousHashes)\n| where FileName =~ \"elzvcf.exe\" and FolderPath has \"AppData\" and FolderPath has \"Roaming\"  // Execution of dropped payload\n| extend ActivityType = \"Process\", EntityType = \"Process\", ActionType = \"ProcessStart\", TargetPath = FolderPath;\n\n// Look for registry changes associated with persistence via the Run key\nlet RegistryEvents = DeviceRegistryEvents\n| where Timestamp > ago(7d)\n| where (RegistryKey has \"Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\" and RegistryValueName has \"elzvcf\")  // Specific persistence for elzvcf\n| extend ActivityType = \"Registry\", EntityType = \"Registry\", ActionType = \"RegistryModify\", TargetPath = RegistryKey;\n\n// Combine file, process, and registry events for a comprehensive view of the attack chain\nunion FileEvents, ProcessEvents, RegistryEvents\n| extend ContentSnippet = extractall(@'Copyright \\(c\\) 2026 Xubuntu', AdditionalFields)  // Tampered TOS indicator (future-dated)\n| where isnotempty(ContentSnippet) or SHA256 in (MaliciousHashes)  // Filter for confirmed malice\n| summarize\n    StartTime = min(Timestamp),\n    EndTime = max(Timestamp),\n    EventCount = count(),\n    ActivityTypes = make_set(ActivityType),\n    EntityTypes = make_set(EntityType),\n    ActionTypes = make_set(ActionType),\n    FilePaths = make_set(TargetPath),\n    FileHashes = make_set(SHA256),\n    RegistryKeys = make_set(RegistryKey),\n    RegistryValues = make_set(RegistryValueName)\n    by DeviceName, DeviceId, InitiatingProcessAccountName\n| extend\n    AlertReason = case(\n        \"Registry\" in (ActivityTypes), \"Persistence via Run key for elzvcf.exe\",\n        \"Process\" in (ActivityTypes), \"Execution of malicious payload elzvcf.exe\",\n        \"File\" in (ActivityTypes) and \"bfaf3af21644ce06fb295c78e080477d9095cb185fcf164ce9cf066a9682b9f1\" in (FileHashes), \"Tampered tos.txt detected\",\n        \"File\" in (ActivityTypes), \"Malicious file associated with Xubuntu watering hole attack detected\",\n        \"Unknown activity detected\"\n    )\n| project DeviceName, DeviceId, StartTime, EndTime, EventCount, ActivityTypes, EntityTypes, ActionTypes, FilePaths, FileHashes, RegistryKeys, RegistryValues, AlertReason\n| order by StartTime desc"
query_type: kql
severity: medium
tactics:
    - TA0002
    - TA0003
    - TA0004
    - TA0001
techniques:
    - T1204.002
    - T1547.001
    - T1189
tags:
    - kql
version: "4"
references:
    - https://cybernews.com/security/xubuntu-site-compromise-hackers-peddle-malware
