# entries/kql/019bdcf6-89c6-72c8-a0f1-0d1507e00712.yaml
id: 019bdcf6-89c6-72c8-a0f1-0d1507e00712
name: Pulsar RAT Detection Query
description: Pulsar RAT is a derivative of Quasar RAT with extensive functionality including keylogging, cryptocurrency wallet clipping, credential theft, file management, remote shell execution, and data exfiltration capabilities. As a modular, open-source remote administration tool designed for Windows systems, Pulsar introduces significant enhancements over its predecessor.
query: "// Pulsar RAT Detection Query\r\n// Detects: UAC Bypass, Persistence, Anti-Analysis, and Network Indicators\r\nlet timeframe = 12h; // Adjust hunting window as needed\r\n\r\n// STAGE 1: Detect Initial UAC Bypass via Temp BAT file and Registry Modification\r\nlet InitialExecution = materialize(\r\n    DeviceProcessEvents\r\n    | where TimeGenerated > ago(timeframe)\r\n    // Look for creation of the suspicious BAT file pattern in Temp\r\n    | where (FolderPath contains @\"\\AppData\\Local\\Temp\\\" and FileName matches regex @\"^[a-f0-9]{8,10}\\.bat$\")\r\n       or (ProcessCommandLine has @\"ms-settings\" and ProcessCommandLine has @\"Shell\\Open\\command\")\r\n    | extend DetectionStage = \"Initial Execution & UAC Bypass\", IndicatorType = \"Persistence Mechanism\"\r\n    | extend ReportID = strcat(\"PULSAR-INIT-\", DeviceId, \"-\", tostring(hash(strcat(FileName, TimeGenerated))))\r\n    | project ReportID, Timestamp = TimeGenerated, DeviceID = DeviceId, DeviceName, AccountName,\r\n              DetectionStage, IndicatorType, FileName, FolderPath,\r\n              AdditionalData1 = ProcessCommandLine, AdditionalData2 = \"\", AdditionalData3 = \"\"\r\n);\r\n\r\n// STAGE 2: Detect Persistence via Scheduled Task\r\nlet TaskPersistence = materialize(\r\n    DeviceProcessEvents\r\n    | where TimeGenerated > ago(timeframe)\r\n    | where ProcessCommandLine has \"schtasks\" \r\n       and (ProcessCommandLine has @\"\\AppData\\Local\\Temp\\\" or ProcessCommandLine contains \"HIGHEST\")\r\n    | extend DetectionStage = \"Persistence\", IndicatorType = \"Scheduled Task Creation\"\r\n    | extend ReportID = strcat(\"PULSAR-TASK-\", DeviceId, \"-\", tostring(hash(strcat(ProcessCommandLine, TimeGenerated))))\r\n    | project ReportID, Timestamp = TimeGenerated, DeviceID = DeviceId, DeviceName, AccountName,\r\n              DetectionStage, IndicatorType, FileName, FolderPath,\r\n              AdditionalData1 = ProcessCommandLine, AdditionalData2 = \"\", AdditionalData3 = \"\"\r\n);\r\n\r\n// STAGE 3: Detect Anti-Virtualization Checks\r\nlet AntiAnalysis = materialize(\r\n    DeviceProcessEvents\r\n    | where TimeGenerated > ago(timeframe)\r\n    | where ProcessCommandLine has_any (\"QEMU\", \"VBOX\", \"VMWARE\", \"VIRTUAL\")\r\n       or (FileName =~ \"cmd.exe\" and ProcessCommandLine contains \"wmic\" and ProcessCommandLine contains \"diskdrive\")\r\n    | extend DetectionStage = \"Defense Evasion\", IndicatorType = \"Anti-Virtualization Check\"\r\n    | extend ReportID = strcat(\"PULSAR-ANTIVM-\", DeviceId, \"-\", tostring(hash(strcat(ProcessCommandLine, TimeGenerated))))\r\n    | project ReportID, Timestamp = TimeGenerated, DeviceID = DeviceId, DeviceName, AccountName,\r\n              DetectionStage, IndicatorType, FileName, FolderPath,\r\n              AdditionalData1 = ProcessCommandLine, AdditionalData2 = \"\", AdditionalData3 = \"\"\r\n);\r\n\r\n// STAGE 4: Detect Suspicious Network Activity\r\nlet NetworkActivity = materialize(\r\n    DeviceNetworkEvents\r\n    | where TimeGenerated > ago(timeframe)\r\n    | where RemotePort >= 1000 and RemotePort <= 65535\r\n    | where InitiatingProcessFileName !in~ (\"chrome.exe\", \"msedge.exe\", \"firefox.exe\", \"svchost.exe\", \"explorer.exe\")\r\n       and (InitiatingProcessFolderPath contains @\"\\AppData\\Local\\Temp\\\" or InitiatingProcessFolderPath contains @\"\\AppData\\Roaming\\\")\r\n    | summarize ConnectionCount = count(), FirstSeen = min(TimeGenerated),\r\n                SampleRemoteIPs = make_set(RemoteIP, 3), SamplePorts = make_set(RemotePort, 3)\r\n              by DeviceID = DeviceId, DeviceName, InitiatingProcessFileName\r\n    | where ConnectionCount > 20 // Beaconing threshold\r\n    | extend DetectionStage = \"Command & Control\", IndicatorType = \"Suspicious Outbound Beaconing\"\r\n    | extend ReportID = strcat(\"PULSAR-NET-\", DeviceID, \"-\", tostring(hash(strcat(InitiatingProcessFileName, FirstSeen))))\r\n    | project ReportID, Timestamp = FirstSeen, DeviceID, DeviceName, AccountName = \"\",\r\n              DetectionStage, IndicatorType, FileName = InitiatingProcessFileName, FolderPath = \"\",\r\n              AdditionalData1 = strcat(\"Connections: \", ConnectionCount, \", Ports: \", SamplePorts),\r\n              AdditionalData2 = strcat(\"IPs: \", SampleRemoteIPs), AdditionalData3 = \"\"\r\n);\r\n\r\n// COMBINE all detection stages into a unified view\r\nlet DetailedEvents = materialize(\r\n    union\r\n        InitialExecution,\r\n        TaskPersistence,\r\n        AntiAnalysis,\r\n        NetworkActivity\r\n    | order by Timestamp desc\r\n);\r\n\r\n// Create summarized alerts for devices showing multiple indicators\r\nDetailedEvents\r\n| summarize\r\n    FirstSeen = min(Timestamp),\r\n    LastSeen = max(Timestamp),\r\n    AllDetectionStages = make_set(DetectionStage),\r\n    IndicatorsFound = count(),\r\n    RelatedReportIDs = make_set(ReportID),\r\n    AffectedAccounts = make_set(AccountName)\r\n    by DeviceID, DeviceName\r\n| extend ThreatScore = case(\r\n    array_length(AllDetectionStages) >= 3, \"HIGH\",\r\n    array_length(AllDetectionStages) == 2, \"MEDIUM\",\r\n    \"LOW\"\r\n)\r\n| where ThreatScore != \"LOW\" // Focus on devices with correlated activity\r\n| extend AlertID = strcat(\"PULSAR-ALERT-\", DeviceID, \"-\", tostring(hash(strcat(DeviceName, FirstSeen))))\r\n| project AlertID, Timestamp = FirstSeen, DeviceID, DeviceName, ThreatScore,\r\n          IndicatorsFound, AllDetectionStages, RelatedReportIDs, AffectedAccounts,\r\n          FirstSeen, LastSeen\r\n| order by ThreatScore desc, IndicatorsFound desc"
query_type: kql
severity: medium
techniques:
    - T1548.002
    - T1497
    - T1053.005
    - T1547.001
    - T1056
    - T1056.001
    - T1059.003
    - T1573.001
tags:
    - kql
version: "2"
references:
    - https://any.run/malware-trends/pulsar/
