# entries/kql/019b2585-b246-720b-981b-8c51c9870203.yaml
id: 019b2585-b246-720b-981b-8c51c9870203
name: Mustang Panda PlugX Delivery via HTML Smuggling
description: This rule detects the multi-stage delivery of PlugX malware, attributed to Mustang Panda, leveraging HTML smuggling, DLL side-loading, and persistence via registry run keys. The detection chain starts with a browser dropping an MSI file (HTML smuggling), followed by msiexec.exe executing a known sideloading executable (gup.exe or robotaskbaricon.exe). Subsequently, a malicious DLL (libcurl.dll or RoboForm.dll) is side-loaded by the executable from the same path, and finally, persistence is established through a registry run key by the same sideloading executable.
query: "// Name: Mustang Panda PlugX Delivery via HTML Smuggling\r\n// Author: Kaung Khant Ko\r\n// Date: 2025-12-16\r\n// Tactic: Initial Access, Execution, Persistence, Defense Evasion\r\n// Technique: T1027.006, T1574.002, T1547.001\r\n// References: https://news.sophos.com/en-us/2025/12/15/game-of-clones-sophos-and-the-mitre-attck-enterprise-2025-evaluations/\r\n// False Positive Sensitivity: Medium\r\n\r\nlet timeframe = 1h;\r\n// Define known executables used for DLL side-loading by Mustang Panda\r\nlet sideload_executables = dynamic([\"gup.exe\", \"robotaskbaricon.exe\"]);\r\n// Define known malicious DLLs used in the side-loading\r\nlet sideload_dlls = dynamic([\"libcurl.dll\", \"RoboForm.dll\"]);\r\n// Define common registry run keys for persistence\r\nlet run_keys = dynamic([\r\n    @\"SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\",\r\n    @\"SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\"\r\n]);\r\n// Define common browser processes that might drop the initial file\r\nlet browser_procs = dynamic([\"msedge.exe\", \"chrome.exe\", \"firefox.exe\", \"iexplore.exe\"]);\r\n\r\n// Step 1: Find MSI files written to disk by a browser, indicating a potential HTML smuggling download.\r\nlet html_smuggling_drop = DeviceFileEvents\r\n| where Timestamp > ago(timeframe)\r\n| where InitiatingProcessFileName in~ (browser_procs)\r\n| where FileName endswith \".msi\"\r\n| project DropTimestamp = Timestamp, DeviceId, MsiFileName = FileName, DroppingProcess = InitiatingProcessFileName;\r\n\r\n// Step 2: Find the execution of the sideloading executable via msiexec.exe.\r\nlet msi_install = DeviceProcessEvents\r\n| where Timestamp > ago(timeframe)\r\n| where FileName in~ (sideload_executables)\r\n| where InitiatingProcessFileName =~ \"msiexec.exe\"\r\n| project InstallTimestamp = Timestamp, DeviceId, SideloadExe = FileName, SideloadExeProcessId = ProcessId, MsiExecCommandLine = InitiatingProcessCommandLine;\r\n\r\n// Step 3: Find the core DLL side-loading event where the DLL is loaded from the same path as the executable.\r\nlet sideloading = DeviceImageLoadEvents\r\n| where Timestamp > ago(timeframe)\r\n| where InitiatingProcessFileName in~ (sideload_executables)\r\n| where FileName in~ (sideload_dlls)\r\n| where isnotempty(InitiatingProcessFolderPath) and InitiatingProcessFolderPath == FolderPath\r\n| project SideloadTimestamp = Timestamp, DeviceId, SideloadExe = InitiatingProcessFileName, SideloadExeProcessId = InitiatingProcessId, LoadedDll = FileName;\r\n\r\n// Step 4: Find persistence created by the same process that performed the side-loading.\r\nlet persistence = DeviceRegistryEvents\r\n| where Timestamp > ago(timeframe)\r\n| where ActionType == \"RegistryValueSet\"\r\n| where InitiatingProcessFileName in~ (sideload_executables)\r\n| where RegistryKey has_any (run_keys)\r\n| project PersistenceTimestamp = Timestamp, DeviceId, SideloadExe = InitiatingProcessFileName, SideloadExeProcessId = InitiatingProcessId, RegistryKeyPath = strcat(RegistryKey, \"\\\\\", RegistryValueName), RegistryValueData;\r\n\r\n// Join the indicators together to form the full attack chain, starting from the MSI installation.\r\nmsi_install\r\n// Correlate the MSI install with the initial browser download.\r\n| join kind=inner (\r\n    html_smuggling_drop\r\n) on DeviceId\r\n| where MsiExecCommandLine has MsiFileName\r\n| where InstallTimestamp between (DropTimestamp .. (DropTimestamp + 5m))\r\n// Correlate with the DLL side-loading event.\r\n| join kind=inner (\r\n    sideloading\r\n) on DeviceId, SideloadExe, SideloadExeProcessId\r\n| where SideloadTimestamp between (InstallTimestamp .. (InstallTimestamp + 5m))\r\n// Correlate with the persistence creation event.\r\n| join kind=inner (\r\n    persistence\r\n) on DeviceId, SideloadExe, SideloadExeProcessId\r\n| where PersistenceTimestamp between (SideloadTimestamp .. (SideloadTimestamp + 5m))\r\n// FP Tuning: Some legitimate installers (e.g., Notepad++) may use gup.exe.\r\n// The combination of the full chain significantly reduces false positives, but if noise occurs,\r\n// consider adding more specific process characteristics or command-line arguments.\r\n| project\r\n    Timestamp = SideloadTimestamp,\r\n    DeviceId,\r\n    DeviceName,\r\n    DroppingProcess,\r\n    MsiFileName,\r\n    MsiExecCommandLine,\r\n    SideloadExe,\r\n    LoadedDll,\r\n    PersistenceRegistryKey = RegistryKeyPath,\r\n    PersistenceValue = RegistryValueData\r\n| extend\r\n    RuleName = \"Mustang Panda PlugX Delivery via HTML Smuggling\",\r\n    Tactic = \"Initial Access, Execution, Persistence, Defense Evasion\",\r\n    Technique = \"T1027.006, T1574.002, T1547.001\"\r\n"
query_type: kql
severity: medium
tactics:
    - TA0005
    - TA0003
    - TA0004
techniques:
    - T1027.006
    - T1574.001
    - T1547.001
tags:
    - kql
author: Kaung Khant Ko
version: "2"
references:
    - https://news.sophos.com/en-us/2025/12/15/game-of-clones-sophos-and-the-mitre-attck-enterprise-2025-evaluations/
