# entries/kql/019cd2e8-cad8-71bd-8420-3ded89ac6ce5.yaml
id: 019cd2e8-cad8-71bd-8420-3ded89ac6ce5
name: Advanced Multi-Stage Linux Enumeration & Post-Exploitation Detector
description: |-
    This query identifies potential post-exploitation behavior on Linux systems by monitoring for clusters of discovery and enumeration commands. Instead of alerting on single, potentially benign commands, it utilizes a behavioral scoring engine that categorizes activities such as credential hunting, privilege escalation, and network scanning.

    By implementing extensive allowlists for common administrative tools (e.g., Ansible, Zabbix, Monitoring Agents) and requiring a minimum threshold of unique activity categories, the query effectively filters out "noise." It calculates a dynamic Risk Score and assigns a Severity level based on the criticality of the commands and the speed of execution (Burst Detection), making it a highly reliable tool for SOC analysts to detect active hands-on-keyboard attacks.
query: "// ============================================================================\r\n// Linux Discovery & Enumeration Detection\r\n// Detect post-exploitation discovery with minimal false positives\r\n// ============================================================================\r\nlet timeframe = 4h;\r\nlet threshold = 6; // Change this threshold value for your needs\r\n// --- ALLOWLISTS (adjust to your environment!) ---\r\n// Known service accounts that routinely perform discovery commands\r\nlet excludedAccounts = dynamic([\r\n    \"zabbix\", \"nagios\", \"icinga\", \"datadog\", \"telegraf\", \"prometheus\",\r\n    \"collectd\", \"splunk\", \"omsagent\", \"waagent\", \"nessus\", \"qualys\",\r\n    \"rapid7\", \"ansible\", \"puppet\", \"chef\", \"salt\", \"cfagent\",\r\n    \"pcp\", \"sysstat\", \"munin\", \"cacti\", \"xymon\"\r\n]);\r\n// Automation and monitoring parent processes\r\nlet excludedParentImages = dynamic([\r\n    \"ansible\", \"ansible-playbook\", \"puppet\", \"chef-client\", \"salt-minion\",\r\n    \"salt-call\", \"cfagent\", \"zabbix_agentd\", \"nagios\", \"nrpe\",\r\n    \"collectd\", \"telegraf\", \"datadog-agent\", \"omsagent\",\r\n    \"omi\", \"auoms\", \"qualys-cloud-agent\", \"nessus\",\r\n    \"cron\", \"crond\", \"anacron\", \"atd\"\r\n]);\r\n// --- PATTERN DEFINITIONS (specific patterns only, no short generic ones) ---\r\nlet suspiciousPatterns = dynamic([\r\n    // User & Identity Discovery (specific, no short forms like \"id \" or \"w \")\r\n    \"whoami\", \"getent passwd\", \"cat /etc/passwd\", \"cat /etc/shadow\", \r\n    \"cat /etc/group\", \"lastlog\",   \r\n    // System Info Discovery\r\n    \"uname -a\", \"cat /etc/os-release\", \"hostnamectl\", \"lsb_release\",\r\n    \"dmidecode\", \"cat /proc/version\", \"cat /proc/cpuinfo\",\r\n    // Privilege Escalation Enumeration (highly suspicious)\r\n    \"sudo -l\", \"find / -perm\", \"find / -suid\", \"find / -sgid\",\r\n    \"find / -writable\", \"getcap -r\",  \r\n    // Persistence Enumeration\r\n    \"crontab -l\", \"cat /etc/crontab\", \"ls /etc/cron\",\r\n    \"systemctl list-unit-files\", \"service --status-all\",\r\n    \"ls /etc/init.d\", \"cat /etc/rc.local\", \r\n    // Process & Service Discovery (specific)\r\n    \"ps aux\", \"ps -ef\", \"pstree\", \"lsof -i\",  \r\n    // Network Discovery\r\n    \"netstat -\", \"ss -tulpn\", \"ss -antp\", \"arp -a\", \"ip addr\", \"ip route\",\r\n    \"ifconfig\", \"route -n\", \"cat /etc/hosts\", \"cat /etc/resolv.conf\",\r\n    \"nmap \", \"masscan \",  \r\n    // Credential Hunting (highly suspicious)\r\n    \"id_rsa\", \"id_dsa\", \"id_ecdsa\", \"id_ed25519\",\r\n    \"cat /etc/ssh\", \".ssh/config\", \".ssh/authorized_keys\",\r\n    \"find / -name *.pem\", \"find / -name *.key\",\r\n    \"find / -name *.pfx\", \"find / -name *.p12\",\r\n    \".bash_history\", \".mysql_history\", \".psql_history\", \r\n    // Environment Discovery (specific)\r\n    \"printenv\", \"cat /etc/environment\", \"cat /proc/net\",    \r\n    // Filesystem / Config Hunting\r\n    \"cat /etc/fstab\",\r\n    \"find / -name wp-config\", \"find / -name config.php\",\r\n    \"find / -name .env\", \"find / -name database.yml\",\r\n    \"find / -name web.config\",   \r\n    // Software / Package Discovery\r\n    \"dpkg -l\", \"rpm -qa\", \"apt list --installed\", \"pip list\", \"pip3 list\",   \r\n    // Container / Cloud Discovery (highly suspicious)\r\n    \"cat /proc/1/cgroup\", \"ls /.dockerenv\", \r\n    \"curl 169.254.169.254\", \"wget 169.254.169.254\",\r\n    \"curl metadata.google\", \"wget metadata.google\",   \r\n    // LOTL Tool Detection\r\n    \"which gcc\", \"which python\", \"which perl\", \"which ruby\",\r\n    \"which nc\", \"which ncat\", \"which socat\",\r\n    \"which wget\", \"which curl\", \"which nmap\"\r\n]);\r\n// --- MAIN QUERY ---\r\nDeviceProcessEvents\r\n| where TimeGenerated > ago(timeframe)\r\n| where ProcessCommandLine has_any (suspiciousPatterns)\r\n// ---- FP FILTER 1: Exclude known service accounts ----\r\n| where AccountName !in (excludedAccounts)\r\n// ---- FP FILTER 2: Exclude known automation parent processes ----\r\n| where InitiatingProcessFileName !in (excludedParentImages)\r\n// ---- FP FILTER 3: Exclude known automation command lines ----\r\n| where InitiatingProcessCommandLine !has \"ansible\"\r\n    and InitiatingProcessCommandLine !has \"puppet\"\r\n    and InitiatingProcessCommandLine !has \"chef\"\r\n    and InitiatingProcessCommandLine !has \"salt-call\"\r\n    and InitiatingProcessCommandLine !has \"cfengine\"\r\n// ---- Categorization with risk weighting ----\r\n| extend NormalizedCommand = case(\r\n    // --- High Risk Categories (Weight 3) ---\r\n    ProcessCommandLine has \"sudo -l\" or ProcessCommandLine has \"find / -perm\" \r\n        or ProcessCommandLine has \"find / -suid\" or ProcessCommandLine has \"find / -sgid\"\r\n        or ProcessCommandLine has \"find / -writable\" or ProcessCommandLine has \"getcap -r\", \r\n        \"Privilege_Escalation_Enum\",\r\n    ProcessCommandLine has \"id_rsa\" or ProcessCommandLine has \"id_dsa\"\r\n        or ProcessCommandLine has \"id_ecdsa\" or ProcessCommandLine has \"id_ed25519\"\r\n        or ProcessCommandLine has \"find / -name *.pem\" or ProcessCommandLine has \"find / -name *.key\"\r\n        or ProcessCommandLine has \"find / -name *.pfx\" or ProcessCommandLine has \".ssh/config\"\r\n        or ProcessCommandLine has \".ssh/authorized_keys\"\r\n        or ProcessCommandLine has \"cat /etc/shadow\", \r\n        \"Credential_Hunting\",\r\n    ProcessCommandLine has \"cat /proc/1/cgroup\" or ProcessCommandLine has \"ls /.dockerenv\"\r\n        or ProcessCommandLine has \"169.254.169.254\" or ProcessCommandLine has \"metadata.google\", \r\n        \"Cloud_Container_Discovery\",\r\n    // --- Medium Risk Categories (Weight 2) ---\r\n    ProcessCommandLine has \"nmap \" or ProcessCommandLine has \"masscan \", \r\n        \"Active_Network_Scanning\",\r\n    ProcessCommandLine has \"crontab -l\" or ProcessCommandLine has \"cat /etc/crontab\"\r\n        or ProcessCommandLine has \"ls /etc/cron\"\r\n        or ProcessCommandLine has \"systemctl list-unit-files\" or ProcessCommandLine has \"service --status-all\"\r\n        or ProcessCommandLine has \"ls /etc/init.d\" or ProcessCommandLine has \"cat /etc/rc.local\", \r\n        \"Persistence_Enum\",\r\n    ProcessCommandLine has \"find / -name wp-config\" or ProcessCommandLine has \"find / -name config.php\"\r\n        or ProcessCommandLine has \"find / -name .env\" or ProcessCommandLine has \"find / -name database.yml\"\r\n        or ProcessCommandLine has \"find / -name web.config\", \r\n        \"Config_File_Hunting\",\r\n    ProcessCommandLine has \".bash_history\" or ProcessCommandLine has \".mysql_history\" \r\n        or ProcessCommandLine has \".psql_history\", \r\n        \"History_Exfil\",\r\n    ProcessCommandLine has \"which gcc\" or ProcessCommandLine has \"which python\"\r\n        or ProcessCommandLine has \"which perl\" or ProcessCommandLine has \"which nc\"\r\n        or ProcessCommandLine has \"which ncat\" or ProcessCommandLine has \"which socat\"\r\n        or ProcessCommandLine has \"which wget\" or ProcessCommandLine has \"which curl\"\r\n        or ProcessCommandLine has \"which nmap\", \r\n        \"LOTL_Tool_Discovery\",\r\n    // --- Low Risk Categories (Weight 1) ---\r\n    ProcessCommandLine has \"whoami\" or ProcessCommandLine has \"getent passwd\"\r\n        or ProcessCommandLine has \"cat /etc/passwd\" or ProcessCommandLine has \"cat /etc/group\"\r\n        or ProcessCommandLine has \"lastlog\", \r\n        \"User_Discovery\",\r\n    ProcessCommandLine has \"uname\" or ProcessCommandLine has \"cat /etc/os-release\" \r\n        or ProcessCommandLine has \"hostnamectl\" or ProcessCommandLine has \"lsb_release\"\r\n        or ProcessCommandLine has \"dmidecode\" or ProcessCommandLine has \"cat /proc/version\", \r\n        \"System_Info_Discovery\",\r\n    ProcessCommandLine has \"ps aux\" or ProcessCommandLine has \"ps -ef\" \r\n        or ProcessCommandLine has \"pstree\" or ProcessCommandLine has \"lsof -i\", \r\n        \"Process_Discovery\",\r\n    ProcessCommandLine has \"netstat\" or ProcessCommandLine has \"ss -tulpn\" or ProcessCommandLine has \"ss -antp\"\r\n        or ProcessCommandLine has \"arp -a\" or ProcessCommandLine has \"ip addr\"\r\n        or ProcessCommandLine has \"ip route\" or ProcessCommandLine has \"ifconfig\"\r\n        or ProcessCommandLine has \"cat /etc/hosts\" or ProcessCommandLine has \"cat /etc/resolv.conf\", \r\n        \"Network_Discovery\",\r\n    ProcessCommandLine has \"printenv\" or ProcessCommandLine has \"cat /etc/environment\"\r\n        or ProcessCommandLine has \"cat /proc/net\", \r\n        \"Environment_Discovery\",\r\n    ProcessCommandLine has \"dpkg -l\" or ProcessCommandLine has \"rpm -qa\"\r\n        or ProcessCommandLine has \"apt list --installed\" or ProcessCommandLine has \"pip list\", \r\n        \"Software_Discovery\",\r\n    \"other\"\r\n)\r\n| where NormalizedCommand != \"other\"\r\n// --- Aggregation per device/user/time window ---\r\n| summarize \r\n    UniqueCategoryCount = dcount(NormalizedCommand), \r\n    DetectedCategories = make_set(NormalizedCommand),\r\n    DetailedCommands = make_set(ProcessCommandLine, 25),\r\n    CommandCount = count(),\r\n    FirstSeen = min(TimeGenerated), \r\n    LastSeen = max(TimeGenerated),\r\n    ReportId = take_any(ReportId) \r\n    by DeviceId, DeviceName, AccountName, bin(TimeGenerated, timeframe)\r\n| where UniqueCategoryCount >= threshold\r\n// --- Risk Scoring: High-risk categories carry more weight ---\r\n| extend RiskScore = \r\n      iif(set_has_element(DetectedCategories, \"Credential_Hunting\"), 3, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Privilege_Escalation_Enum\"), 3, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Cloud_Container_Discovery\"), 3, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Active_Network_Scanning\"), 2, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Persistence_Enum\"), 2, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Config_File_Hunting\"), 2, 0)\r\n    + iif(set_has_element(DetectedCategories, \"History_Exfil\"), 2, 0)\r\n    + iif(set_has_element(DetectedCategories, \"LOTL_Tool_Discovery\"), 2, 0)\r\n    + iif(set_has_element(DetectedCategories, \"User_Discovery\"), 1, 0)\r\n    + iif(set_has_element(DetectedCategories, \"System_Info_Discovery\"), 1, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Process_Discovery\"), 1, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Network_Discovery\"), 1, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Environment_Discovery\"), 1, 0)\r\n    + iif(set_has_element(DetectedCategories, \"Software_Discovery\"), 1, 0)\r\n// --- Burst Detection: Timespan of activity ---\r\n| extend ActivityDurationMinutes = datetime_diff('minute', LastSeen, FirstSeen)\r\n// --- Severity Classification ---\r\n| extend Severity = case(\r\n    RiskScore >= 12 and ActivityDurationMinutes <= 30, \"Critical\",\r\n    RiskScore >= 8, \"High\",\r\n    RiskScore >= 5 and UniqueCategoryCount >= 7, \"Medium\",\r\n    \"Low\"\r\n)\r\n// --- Only output relevant alerts ---\r\n| where Severity in (\"Critical\", \"High\", \"Medium\")\r\n// --- Final Projection ---\r\n| project \r\n    Timestamp=TimeGenerated,\r\n    DeviceId,\r\n    DeviceName, \r\n    AccountName, \r\n    Severity,\r\n    RiskScore,\r\n    UniqueCategoryCount, \r\n    CommandCount,\r\n    ActivityDurationMinutes,\r\n    DetectedCategories, \r\n    DetailedCommands,\r\n    ReportId\r\n| order by RiskScore desc, UniqueCategoryCount desc"
query_type: kql
severity: medium
tactics:
    - TA0007
    - TA0004
    - TA0006
    - TA0003
    - TA0002
techniques:
    - T1518
    - T1087
    - T1068
    - T1053.003
    - T1569.003
    - T1083
    - T1057
    - T1049
    - T1016
    - T1135
    - T1552.005
    - T1611
    - T1033
    - T1003.008
    - T1552.003
    - T1059.006
    - T1037.004
tags:
    - kql
author: Benjamin Zulliger
version: "1"
