# entries/kql/019ca507-9e6b-70a4-87c0-34790549dec8.yaml
id: 019ca507-9e6b-70a4-87c0-34790549dec8
name: Moonrise RAT - Go Binary Persistence via Run Key or Startup Folder
description: This rule detects when a process, potentially a Go-compiled binary (indicated by 'go' in the command line), establishes persistence by creating or modifying a registry run key or by dropping a file in a startup folder. This behavior is associated with the Moonrise RAT.
query: "// Name: Moonrise RAT - Go Binary Persistence\r\n// Author: Ethan Andrews + detections.ai\r\n// Date: 2026-02-28\r\n// Description: This rule detects when a process, potentially a Go-compiled binary, establishes persistence by creating or modifying a registry run key or by dropping a file in a startup folder. This behavior is associated with the Moonrise RAT.\r\n// Tactic: Persistence\r\n// Technique: T1547.001\r\n// False Positive Sensitivity: Medium\r\n// References: https://any.run/cybersecurity-blog/moonrise-rat-detected/#indicatorsofcompromiseiocs-18719\r\nlet timeframe = 1h;\r\n// Define common persistence registry run keys\r\nlet runKeys = dynamic([\r\n    \"\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\",\r\n    \"\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\",\r\n    \"\\\\Software\\\\Wow6432Node\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\",\r\n    \"\\\\Software\\\\Wow6432Node\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\"\r\n]);\r\n// Define common persistence startup folders\r\nlet startupFolders = dynamic([\r\n    \"\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\\"\r\n]);\r\n// Part 1: Detect registry persistence\r\nlet registryPersistence =\r\nDeviceRegistryEvents\r\n| where Timestamp > ago(timeframe)\r\n| where ActionType in (\"RegistryKeySet\", \"RegistryValueSet\")\r\n| where RegistryKey has_any (runKeys)\r\n// The detection logic identifies Go-lang binaries by looking for \"go\" in the command line.\r\n// This is a weak heuristic and may require tuning or replacement with more specific indicators.\r\n| where InitiatingProcessCommandLine has \"go\"\r\n| extend PersistenceMethod = \"Registry\", Details = strcat(RegistryKey, \": \", RegistryValueData)\r\n| project Timestamp, DeviceId, DeviceName, InitiatingProcessFileName, InitiatingProcessFolderPath, InitiatingProcessCommandLine, PersistenceMethod, Details;\r\n// Part 2: Detect startup folder persistence\r\nlet startupPersistence =\r\nDeviceFileEvents\r\n| where Timestamp > ago(timeframe)\r\n| where ActionType in (\"FileCreated\", \"FileRenamed\")\r\n| where FolderPath has_any (startupFolders) and (FileName) has_any (\"exe\", \"lnk\", \"bat\", \"vbs\")\r\n// The detection logic identifies Go-lang binaries by looking for \"go\" in the command line.\r\n// This is a weak heuristic and may require tuning or replacement with more specific indicators.\r\n| where InitiatingProcessCommandLine has \"go\"\r\n| extend PersistenceMethod = \"Startup Folder\", Details = tostring(FolderPath)\r\n| project Timestamp, DeviceId, DeviceName, InitiatingProcessFileName, InitiatingProcessFolderPath, InitiatingProcessCommandLine, PersistenceMethod, Details;\r\n// Union the results from both persistence methods\r\n(union registryPersistence, startupPersistence)\r\n// FP Tuning: Exclude common software locations and known legitimate updaters/installers.\r\n// This exclusion list should be customized for your specific environment.\r\n| where InitiatingProcessFolderPath !startswith @\"C:\\Program Files\"\r\n    and InitiatingProcessFolderPath !startswith @\"C:\\Windows\\System32\"\r\n| where InitiatingProcessFileName !in~ ('msiexec.exe', 'setup.exe', 'update.exe', 'installer.exe', 'GoogleUpdate.exe', 'OneDrive.exe')\r\n"
query_type: kql
severity: medium
tactics:
    - TA0003
    - TA0004
techniques:
    - T1547.001
tags:
    - kql
author: Ethan Andrews
version: "2"
references:
    - https://any.run/cybersecurity-blog/moonrise-rat-detected/#indicatorsofcompromiseiocs-18719
