# entries/kql/019d0c27-6e44-729b-8432-a7bd15447085.yaml
id: 019d0c27-6e44-729b-8432-a7bd15447085
name: TernDoor Persistence via Hidden Scheduled Task
description: Detects the creation of the 'WSPrint' scheduled task followed by specific registry modifications intended to hide the task from the Task Scheduler UI. This is a known persistence technique used by the UAT-9244 actor for its TernDoor backdoor.
query: "// Name: UAT-9244 TernDoor Persistence via Hidden Scheduled Task\r\n// Author:Brandon + detections.ai\r\n// Date: 2026-03-20\r\n// Description: Detects the creation of the 'WSPrint' scheduled task followed by specific registry modifications intended to hide the task from the Task Scheduler UI. This is a known persistence technique used by the UAT-9244 actor for its TernDoor backdoor.\r\n// References: https://blog.talosintelligence.com/uat-9244/\r\n// MITRE ATT&CK: T1053.005, T1562.001, T1112\r\nlet timeFrame = 30d;\r\nlet taskName = \"WSPrint\";\r\nlet regPathFragment = strcat(@\"\\Schedule\\TaskCache\\Tree\\\", taskName);\r\n// Collect process events for the specific scheduled task creation\r\nlet taskCreationEvents =\r\n    DeviceProcessEvents\r\n    | where TimeGenerated > ago(timeFrame)\r\n    | project-rename taskCreationEvents_TimeGenerated = TimeGenerated\r\n    | project-rename taskCreationEvents_ProcessCommandLine = ProcessCommandLine\r\n    | where FileName =~ \"schtasks.exe\"\r\n    | where taskCreationEvents_ProcessCommandLine has_all (\"/create\", \"/tn\", taskName);\r\n// Collect registry events related to hiding the task from the UI\r\nlet registryHidingEvents =\r\n    DeviceRegistryEvents\r\n    | where TimeGenerated > ago(timeFrame)\r\n    | project-rename registryHidingEvents_TimeGenerated = TimeGenerated\r\n    | where RegistryKey endswith regPathFragment\r\n    | project-rename registryHidingEvents_RegistryKey = RegistryKey\r\n    | project-rename registryHidingEvents_ActionType = ActionType\r\n    | project-rename registryHidingEvents_RegistryValueName = RegistryValueName\r\n    | project-rename registryHidingEvents_RegistryValueData = RegistryValueData\r\n    // The actor either deletes the Security Descriptor (SD) value or sets the Index value to 0 to hide the task\r\n    | where (registryHidingEvents_ActionType =~ \"RegistryValueDeleted\" and registryHidingEvents_RegistryValueName =~ \"SD\") or (registryHidingEvents_ActionType =~ \"RegistryValueSet\" and registryHidingEvents_RegistryValueName =~ \"Index\" and registryHidingEvents_RegistryValueData == \"0\");\r\n// Join the task creation and registry hiding events on the same device\r\ntaskCreationEvents\r\n| join kind=inner (\r\n    registryHidingEvents\r\n) on DeviceId\r\n// Correlate events where the registry modification occurs within 5 minutes after the task is created to ensure they are related\r\n| where registryHidingEvents_TimeGenerated between (taskCreationEvents_TimeGenerated .. (taskCreationEvents_TimeGenerated + 5m))\r\n// Group events to produce a single alert per incident\r\n| summarize\r\n    StartTime = min(taskCreationEvents_TimeGenerated),\r\n    EndTime = max(registryHidingEvents_TimeGenerated),\r\n    TaskCreationCommandLine = take_any(taskCreationEvents_ProcessCommandLine),\r\n    RegistryKey = take_any(registryHidingEvents_RegistryKey),\r\n    RegistryActions = make_set(pack(\"Action\", registryHidingEvents_ActionType, \"ValueName\", registryHidingEvents_RegistryValueName, \"ValueData\", registryHidingEvents_RegistryValueData))\r\n    by DeviceId, DeviceName, InitiatingProcessAccountName\r\n// FP Tuning Suggestion: If legitimate software performs similar actions, you may need to exclude specific InitiatingProcessFileName values or command lines.\r\n| extend\r\n    timestamp = StartTime,\r\n    AccountCustomEntity = InitiatingProcessAccountName,\r\n    HostCustomEntity = DeviceName\r\n"
query_type: kql
severity: medium
tactics:
    - TA0002
    - TA0003
    - TA0004
    - TA0005
techniques:
    - T1053.005
    - T1112
tags:
    - kql
version: "2"
references:
    - https://blog.talosintelligence.com/uat-9244/
