# entries/kql/5916ab37-583c-44b4-aa7c-3bdcd1ead7ba.yaml
id: 5916ab37-583c-44b4-aa7c-3bdcd1ead7ba
name: WinRAR Zero-Day  ( CVE-2025-8088)
description: KQL Query designed to detected tactics employed during active exploitation of WinRAR Zero-Day  ( CVE-2025-8088)
query: "\nDeviceFileEvents\n| where InitiatingProcessFileName has_any(\"Winrar\", \"Rar\", \"unrar.exe\")\n| where ActionType == \"FileCreated\" \n| where FolderPath contains @\"AppData\\Local\\Temp\" and (FileName endswith \".exe\" or FileName endswith \".dll\") or\n        FolderPath contains @\"Start Menu\\Programs\\Startup\" and (FileName endswith \".lnk\")\n|project ReportId, Timestamp, DeviceName, DeviceId, FileName, InitiatingProcessCommandLine"
query_type: kql
severity: medium
tactics:
    - TA0009
    - TA0003
    - TA0004
techniques:
    - T1560.001
    - T1547.001
tags:
    - kql
version: "1"
