# entries/kql/9406ce79-15c7-40b7-b2a3-b8e3d24f7d40.yaml
id: 9406ce79-15c7-40b7-b2a3-b8e3d24f7d40
name: Hunt for Svchost or Schtasks Initiated Processes on Specific Devices
description: This rule detects process events where `svchost.exe` or `schtasks.exe` are the initiating processes on specific devices. This is designed to identify unexpected or suspicious scheduled tasks, whether system-initiated or user-triggered, on monitored devices. It can be used for hunting or to investigate alerts from other security tools regarding suspicious scheduled tasks.
query: "#### This query analyzes process events on specific devices in your network to identify actions involving svchost.exe or schtasks.exe. These executables are commonly associated with scheduled tasks, either system-initiated or user-triggered. \n#### This query can be helpfull in a situation where you want to see or hunt for any unexpected or suspicious scheduled task running automatically on the devices that you want to check. This query can also help in a situation where defender detected a suspicious scheduled task running on the organizations devices.\n\n#### Query:\n```KQL\nDeviceProcessEvents\n| where InitiatingProcessFileName == \"svchost.exe\" // Common for scheduled tasks \n or InitiatingProcessFileName == \"schtasks.exe\" //If the task is directly triggered by the user\n| where DeviceName in (\"adc\", \"sgf\", \"ert\", \"poi\")\n| order by Timestamp desc\n```\n"
query_type: kql
tactics:
    - TA0002
    - TA0003
    - TA0004
techniques:
    - T1053
    - T1053.005
tags:
    - kql
version: "1"
