# entries/kql/a58bd092-57f7-476a-ab27-5f853b4a4b40.yaml
id: a58bd092-57f7-476a-ab27-5f853b4a4b40
name: Startup Registry Key Creation or Modification
description: |-
    This detection identifies attempts to establish persistence by modifying known Windows startup registry keys. These keys control which programs are automatically executed during system boot or user logon and are frequently abused by adversaries to maintain access across reboots.

    Analysts will need to tune the rule through testing and using the exclusion statements to remove noise of what is known good within their environments.
query: "// List of startup registry keys to monitor\r\nlet startupRegistryList = dynamic([\r\n  // HKEY_CURRENT_USER - User-specific autostarts\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunServices',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunServicesOnce',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Policies\\\\Explorer\\\\Run',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\Shell Folders',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\User Shell Folders',\r\n  'HKEY_CURRENT_USER\\\\SOFTWARE\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Windows',\r\n  // HKEY_LOCAL_MACHINE - Machine-wide autostarts\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnceEx',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunServices',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunServicesOnce',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Policies\\\\Explorer\\\\Run',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\Shell Folders',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\User Shell Folders',\r\n  // HKEY_LOCAL_MACHINE - Winlogon-related autostarts\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Winlogon\\\\Shell',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Winlogon\\\\Userinit',\r\n  // HKEY_LOCAL_MACHINE - Boot-level autostarts\r\n  'HKEY_LOCAL_MACHINE\\\\SYSTEM\\\\CurrentControlSet\\\\Control\\\\Session Manager',\r\n  // Wow6432Node - 32-bit app equivalents\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Wow6432Node\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Wow6432Node\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce',\r\n  'HKEY_LOCAL_MACHINE\\\\SOFTWARE\\\\Wow6432Node\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnceEx'\r\n]);\r\n//InitiatingProcessCommandLine Exclusion List \r\nlet exclusionList_IPCL = dynamic([\r\n    'ProcessCommandLineExclusions'\r\n  ]);\r\n//InitiatingProcessName Exclusion List\r\nlet exclusionList_IPFN = dynamic([\r\n    'ProcessNameExclusions',\r\n  ]);\r\n//RegistryValueData Exclusion List\r\nlet exclusionList_RVD = dynamic([\r\n    'RegistryValueExclusions'\r\n  ]);\r\n //InitiatingProcessSHA256 Exclusion List\r\nlet exclusionList_IPSHA256 = dynamic([\r\n    'SHA256HashExclusions'\r\n  ]);\r\nDeviceRegistryEvents\r\n| where ActionType in ('RegistryValueSet', 'RegistryKeyCreated') and tostring(RegistryKey) has_any (startupRegistryList)\r\n| where not (tostring(InitiatingProcessCommandLine) has_any (exclusionList_IPCL))\r\n| where not (tostring(InitiatingProcessFileName) has_any(exclusionList_IPFN))\r\n| where not (tostring(RegistryValueData) has_any(exclusionList_RVD))\r\n| where not (tostring(InitiatingProcessSHA256) has_any(exclusionList_IPSHA256))"
query_type: kql
severity: medium
tactics:
    - TA0003
techniques:
    - T1547.001
tags:
    - kql
version: "1"
