# entries/kql/a91aeeed-71ab-4df3-ae66-b45d0d15a5d2.yaml
id: a91aeeed-71ab-4df3-ae66-b45d0d15a5d2
name: CVE-2025-8088 WinRAR 0-day Path Traversal Exploited via COM Hijacking
description: Detects suspicious activities related to CVE-2025-8088 WinRAR exploitation, including malicious file creation, COM hijacking, C2 communication, and anti-analysis techniques. Enhanced with ADS detection, command-line analysis, and refined false positive handling. This threat is now attributed to RomCom.
query: "// Name: CVE-2025-8088 Winrar <=7.13 0-day Path Traversal Exploitated via COM Hijacking \r\n// Author: Reyben T. Cortes @ crisac.org\r\n// Date: 2025-08-13\r\n// Description: Detects suspicious activities related to CVE-2025-8088 WinRAR exploitation, including malicious file creation, COM hijacking, C2 communication, and anti-analysis techniques. Enhanced with ADS detection, command-line analysis, and refined false positive handling. This threat is now attributed to RomCom.\r\n// Reference: https://www.welivesecurity.com/en/eset-research/update-winrar-tools-now-romcom-and-others-exploiting-zero-day-vulnerability/\r\n// Techniques: T1546.015, T1204.002, T1071.001, T1497\r\n// False Positive Sensitivity: Low\r\n// Note: Legitimate software installers may write unsigned DLLs to AppData. Tune by adding known good publisher names or specific command-line arguments.\r\n\r\nlet SuspiciousCLSID = \"{1299CF18-C4F5-4B6A-BB0F-2299F0398E27}\";\r\nlet PathTraversalPatterns = dynamic([\".../...//\", \"../\", \"..\\\\\", \"%2e%2e/\", \"%2e%2e%5c\", \"..\\\\..\\\\\", \":\\\\\", \"%252e%252e/\", \":\"]);\r\nlet KnownMaliciousSHA1 = dynamic([\r\n    \"AE687BEF963CB30A3788E34CC18046F54C41FFBA\", // msedge.dll\r\n    \"F77DBA76010A9988C9CEB8E420C96AEBC071B889\", // Eli_Rosenfeld_CV2 - Copy (10).rar\r\n    \"371A5B8BA86FBCAB80D4E0087D2AA0D8FFDDC70B\", // Adverse_Effect_Medical_Records_2025.rar\r\n    \"D43F49E6A586658B5422EDC647075FFD405D6741\", // cv_submission.rar\r\n    \"C340625C779911165E3983C77FD60855A2575275\", // cv_submission.rar\r\n    \"676086860055F6591FED303B4799C725F8466CF4\", // Datos adjuntos sin título 00170.dat\r\n    \"1F25E062E8E9A4F1792C3EAC6462694410F0F1CA\", // JobDocs_July2025.rar\r\n    \"C94A6BD6EC88385E4E831B208FED2FA6FAED6666\", // Recruitment_Dossier_July_2025.rar\r\n    \"01D32FE88ECDEA2B934A00805E138034BF85BF83\", // install_module_x64.dll\r\n    \"AB79081D0E26EA278D3D45DA247335A545D0512E\", // Complaint.exe\r\n    \"1AEA26A2E2A7711F89D06165E676E11769E2FD68\"  // ApbxHelper.exe\r\n]);\r\nlet SuspiciousDomains = dynamic([\"srlaptop.com\", \"campanole.com\", \"melamorri.com\", \"gohazeldale.com\"]);\r\nlet SuspiciousIPs = dynamic([\"162.19.175.44\", \"194.36.209.127\", \"85.158.108.62\", \"185.173.235.134\"]);\r\nlet SuspiciousPaths = dynamic([\"AppData\\\\Local\\\\Temp\", \"AppData\\\\Local\", \"Start Menu\\\\Programs\\\\Startup\"]);\r\nlet TrustedPublishers = dynamic([\"Microsoft Corporation\", \"WinRAR\", \"RARLAB\"]);\r\nlet TimeWindow = ago(30d);\r\n\r\n// Registry events for COM hijacking and path traversal\r\nlet RegistryEvents = DeviceRegistryEvents\r\n| where Timestamp > TimeWindow\r\n| where ActionType == \"RegistryValueSet\"\r\n| where (RegistryKey has \"CLSID\\\\\" and RegistryKey has \"InprocServer32\" and (RegistryValueData has_any (PathTraversalPatterns) or RegistryValueData has_any (SuspiciousPaths)))\r\n   or (RegistryKey has \"CLSID\\\\\" + SuspiciousCLSID + \"\\\\InprocServer32\" and RegistryValueData has \"msedge.dll\")\r\n| project EventTime = Timestamp, DeviceId, DeviceName, InitiatingProcessFileName, InitiatingProcessCommandLine, RegistryKey, RegistryValueData, ReportId, EventType=\"Registry\";\r\n\r\n// File events for malicious file deployment by WinRAR\r\nlet FileEvents = DeviceFileEvents\r\n| where Timestamp > TimeWindow\r\n| where ActionType in (\"FileCreated\", \"FileModified\")\r\n| where (FileName endswith_any (\"msedge.dll\", \"ApbxHelper.exe\", \"Complaint.exe\", \"install_module_x64.dll\", \".lnk\", \"cv_submission.rar\", \"JobDocs_July2025.rar\", \"Recruitment_Dossier_July_2025.rar\", \"Datos adjuntos sin título 00170.dat\", \"Adverse_Effect_Medical_Records_2025.rar\")\r\n        or FileName matches regex @\".\\.(dll|exe|rar|dat|cmd|ps1|bat|dat|lnk)$\"\r\n        or FileName startswith \"Eli_Rosenfeld_CV2\")\r\n| where FolderPath has_any (SuspiciousPaths)\r\n// The initiating process is key to this specific threat.\r\n| where InitiatingProcessFileName == \"WinRAR.exe\"\r\n// Check for WinRAR extraction commands to reduce FPs from just viewing archive contents.\r\n| where InitiatingProcessCommandLine matches regex @\"(x|e)\\s+.*(\\.rar|\\.dat)\"\r\n// Look for known malicious files or any unsigned file not from a trusted publisher.\r\n| where (SHA1 in (KnownMaliciousSHA1) or (IsSigned == false and Signer !in (TrustedPublishers)))\r\n| project EventTime = Timestamp, DeviceId, DeviceName, InitiatingProcessFileName, InitiatingProcessCommandLine, FileName, FolderPath, SHA1, Signer, ReportId, EventType=\"File\";\r\n\r\n// Network events for C2 communication\r\nlet NetworkEvents = DeviceNetworkEvents\r\n| where Timestamp > TimeWindow\r\n| where RemoteUrl has_any (SuspiciousDomains) or RemoteIP in (SuspiciousIPs) or RemoteUrl matches regex @\"/s/0\\.7\\.8/clarity\\.js|/TOfrPOseJKZ|/iEZGPctehTZ\"\r\n| project EventTime = Timestamp, DeviceId, DeviceName, InitiatingProcessFileName, RemoteUrl, RemoteIP, ReportId, EventType=\"Network\";\r\n\r\n// Process events that may indicate execution of the malicious payload\r\nlet ProcessEvents = DeviceProcessEvents\r\n| where Timestamp > TimeWindow\r\n| where FileName in (\"msedge.dll\", \"ApbxHelper.exe\", \"Complaint.exe\")\r\n| where ProcessCommandLine has_any (SuspiciousDomains) // Check for hardcoded domains in command line\r\n| project EventTime = Timestamp, DeviceId, DeviceName, InitiatingProcessFileName = FileName, InitiatingProcessCommandLine = ProcessCommandLine, ReportId, EventType=\"Process\";\r\n\r\n// Enrichment data: check for RecentDocs flooding on a device\r\nlet RecentDocsFlooding = DeviceRegistryEvents\r\n| where Timestamp > TimeWindow\r\n| where RegistryKey has \"SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\RecentDocs\"\r\n| summarize makeset(RegistryValue), DocCount = count() by DeviceId\r\n| where DocCount >= 69;\r\n\r\n// Combine all primary event streams\r\nlet AllEvents = union RegistryEvents, FileEvents, NetworkEvents, ProcessEvents;\r\n\r\nAllEvents\r\n// Enrich events with information about RecentDocs flooding on the same device\r\n| join kind=leftouter (\r\n    RecentDocsFlooding\r\n) on DeviceId\r\n// Enrich events with a summary of all suspicious files dropped on the same device\r\n| join kind=leftouter (\r\n    FileEvents\r\n    | summarize FileDetails = make_set(pack(\"File\", FileName, \"Path\", FolderPath, \"SHA1\", SHA1, \"Signer\", Signer)) by DeviceId\r\n) on DeviceId\r\n| project-reorder EventTime, DeviceId, DeviceName, EventType, InitiatingProcessFileName, InitiatingProcessCommandLine\r\n| project-away DeviceId1, DeviceId2\r\n| order by EventTime desc\r\n"
query_type: kql
severity: medium
tactics:
    - TA0004
    - TA0003
    - TA0002
    - TA0011
    - TA0005
    - TA0007
techniques:
    - T1546.015
    - T1204.002
    - T1071.001
    - T1497
tags:
    - kql
author: Reyben Cortes @ crisac.org
version: "1"
references:
    - https://www.welivesecurity.com/en/eset-research/update-winrar-tools-now-romcom-and-others-exploiting-zero-day-vulnerability/
