# entries/sigma/019a69c1-3bf7-7140-b980-bc8007ef1688.yaml
id: 019a69c1-3bf7-7140-b980-bc8007ef1688
name: LeakyStealer Registry Persistence
description: Detects the creation of the "EdgeUpdateCore" Run key, a persistence mechanism used by the LeakyStealer malware.
query: "title: LeakyStealer Registry Persistence\r\nid: 1c4f11bb-ffde-451b-a559-789feacd8d83\r\nstatus: experimental\r\ndescription: Detects the creation of the \"EdgeUpdateCore\" Run key, a persistence\r\n  mechanism used by the LeakyStealer malware.\r\nreferences:\r\n  - https://hybrid-analysis.blogspot.com/2025/11/leakyinjector-and-leakystealer-duo.html\r\nauthor: detections.ai, Georg Lauenstein\r\ndate: 2025-11-09\r\ntags:\r\n  - attack.t1547.001\r\nlogsource:\r\n  category: registry_set\r\n  product: windows\r\ndetection:\r\n  selection:\r\n    # Detects the specific Run key value name used by LeakyStealer\r\n    TargetObject|contains: '\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\EdgeUpdateCore' # Detects the specific executable name and path in the AppData directory\r\n    Details|contains|all:\r\n      - 'AppData'\r\n      - 'MicrosoftEdgeUpdateCore.exe'\r\n  condition: selection\r\nfalsepositives:\r\n  - Legitimate Microsoft Edge update processes.\r\nlevel: high\r\n"
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0004
techniques:
    - T1547.001
tags:
    - attack.t1547.001
author: detections.ai, Georg Lauenstein
version: "3"
references:
    - https://hybrid-analysis.blogspot.com/2025/11/leakyinjector-and-leakystealer-duo.html
data_sources:
    - registry_set
    - windows
