# entries/sigma/019abed0-34e3-755c-a5cc-58feafb0401c.yaml
id: 019abed0-34e3-755c-a5cc-58feafb0401c
name: FlexibleFerret LaunchAgent Persistence Detection
description: Detects creation of FlexibleFerret persistence mechanism via LaunchAgent plist in macOS.
query: |+
    ---
    title: FlexibleFerret LaunchAgent Persistence
    id: a1b2c3d4-4444-4e5f-8a9b-0c1d2e3f4a5b
    status: experimental
    description: Detects creation of FlexibleFerret persistence mechanism via LaunchAgent plist
    author: Claude Opus 4.5
    date: 2025/11/26
    references:
        - https://www.jamf.com/blog/flexibleferret-malware-continues-to-adapt/
    tags:
        - attack.persistence
        - attack.t1543.001
        - attack.t1547.011
    logsource:
        category: file_event
        product: macos
    detection:
        selection_path:
            TargetFilename|contains: 'Library/LaunchAgents/'
        selection_filename:
            TargetFilename|contains:
                - 'com.driver9990as7tpatch.plist'
                - 'driver'
                - 'patch.plist'
        condition: selection_path and selection_filename
    falsepositives:
        - Legitimate driver-related LaunchAgents (review naming patterns)
    level: high

query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0004
techniques:
    - T1543.001
tags:
    - attack.persistence
    - attack.t1543.001
    - attack.t1547.011
author: Claude Opus 4.5
version: "1"
references:
    - https://www.jamf.com/blog/flexibleferret-malware-continues-to-adapt/
data_sources:
    - file_event
    - macos
