# entries/sigma/019be644-42e0-70df-b547-1a93bb75eccc.yaml
id: 019be644-42e0-70df-b547-1a93bb75eccc
name: TamperedChef Registry Run Key Persistence
description: Detects the creation of a specific registry run key used by the TamperedChef malware for persistence.
query: |
    title: TamperedChef Registry Run Key Persistence
    id: 5f8f2472-7b60-4aed-a64f-8d1e93ac8999
    status: experimental
    description: Detects the creation of a specific registry run key used by the
      TamperedChef malware for persistence.
    references:
      - https://www.sophos.com/en-us/blog/tamperedchef-serves-bad-ads-with-infostealers-as-the-main-course
    author: detections.ai
    date: 2026-01-22
    tags:
      - attack.persistence
      - attack.t1547.001
    logsource:
      category: registry_set
      product: windows
    detection:
      selection:
        # Detects the specific registry run key created by TamperedChef
        TargetObject|endswith: '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\PDFEditorUpdater'
      condition: selection
    falsepositives:
      - This detection is highly specific to the TamperedChef malware. False
        positives are not expected.
    level: high
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0004
techniques:
    - T1547.001
tags:
    - attack.persistence
    - attack.t1547.001
author: detections.ai
version: "3"
references:
    - https://www.sophos.com/en-us/blog/tamperedchef-serves-bad-ads-with-infostealers-as-the-main-course
data_sources:
    - registry_set
    - windows
