# entries/sigma/019b9e52-b12d-7218-a4f1-a927106909d8.yaml
id: 019b9e52-b12d-7218-a4f1-a927106909d8
name: Scheduled Task Created with Elevated Privilege Parameters
description: Detects scheduled tasks created on remote systems using explicit credentials.
query: "title: Scheduled Task Created with Elevated Privilege Parameters\r\nid: 7b1350be-03bc-4dd5-886a-53c5a5dafceb\r\nstatus: test\r\ndescription: Detects scheduled tasks created on remote systems using explicit credentials.\r\nreferences:\r\n    - https://attack.mitre.org/techniques/T1053/005/\r\n    - https://attack.mitre.org/techniques/T1548/\r\nauthor: SimSpace\r\ndate: 2026-01-01\r\ntags:\r\n    - attack.privilege-escalation\r\n    - attack.t1053.005\r\n    - attack.t1548\r\nlogsource:\r\n    category: process_creation\r\n    product: windows\r\ndetection:\r\n    selection_schtasks:\r\n        Image|endswith: '\\schtasks.exe'\r\n        CommandLine|contains|all:\r\n            - '/create'\r\n            - '/S '\r\n            - '/U '\r\n            - '/P '\r\n    condition: selection_schtasks\r\nfalsepositives:\r\n    - Unknown\r\nlevel: high"
query_type: sigma
severity: medium
tactics:
    - TA0004
    - TA0002
    - TA0003
    - TA0005
techniques:
    - T1053.005
    - T1548
tags:
    - attack.privilege-escalation
    - attack.t1053.005
    - attack.t1548
author: SimSpace
version: "2"
references:
    - https://attack.mitre.org/techniques/T1053/005/
    - https://attack.mitre.org/techniques/T1548/
data_sources:
    - process_creation
    - windows
