# entries/sigma/019bebfd-e173-73fa-a4ea-b85e39ae53f0.yaml
id: 019bebfd-e173-73fa-a4ea-b85e39ae53f0
name: Fake PDF Converter Persistence Task
description: Detects the creation of scheduled tasks used for persistence by fake PDF converter malware. These tasks are created in the registry under the TaskCache tree.
query: |
    title: Fake PDF Converter Persistence Task
    id: 3432300e-9c18-480e-96f5-8bb3aa69d2f8
    status: experimental
    description: Detects the creation of scheduled tasks used for persistence by
      fake PDF converter malware. These tasks are created in the registry under the
      TaskCache tree.
    references:
      - Internal ID: 78b624b602dd7221c89da3f71acf4500db4e0b2e6f3f00adb894128e22f37695
    author: detections.ai
    date: 2026-01-23
    tags:
      - attack.persistence
      - attack.t1053.005
    logsource:
      category: registry_set
      product: windows
    detection:
      selection:
        # Detects registry keys created for specific scheduled tasks
        TargetObject|contains: '\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tree\'
        TargetObject|endswith:
          - '\ConvertMateTask'
          - '\Crystal_updater'
          - '\Easy2ConvertTask'
          - '\PDC_Update'
      condition: selection
    falsepositives:
      - Unlikely
    level: high
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0002
    - TA0004
techniques:
    - T1053.005
tags:
    - attack.persistence
    - attack.t1053.005
author: detections.ai
version: "2"
references:
    - '[object Object]'
data_sources:
    - registry_set
    - windows
