# entries/sigma/019c48f6-31e2-7267-9470-e82ccead5beb.yaml
id: 019c48f6-31e2-7267-9470-e82ccead5beb
name: Remote Scheduled Task Deletion After Creation
description: Detects the deletion of scheduled tasks on remote systems shortly after their creation, indicating potential cleanup of lateral movement artifacts. This behavior is often observed when attackers create tasks for execution and then remove them to evade detection. The correlation groups events by user, meaning it can match operations performed by the same account across different target systems. This rule specifically looks for `schtasks.exe` being used with the `/create` and `/delete` commands targeting remote systems (excluding localhost, 127.0.0.1, and ::1) within a one-hour window.
query: "title: Remote Scheduled Task Creation\r\nid: 0f0dc15f-1ac9-4982-9ec0-5b813ac89033\r\nname: remote_schtask_creation\r\nlogsource:\r\n    product: windows\r\n    category: process_creation\r\ndetection:\r\n    selection_image:\r\n        Image|endswith: '\\schtasks.exe'\r\n    selection_create:\r\n        CommandLine|contains|all:\r\n            - '/create'\r\n            - '/S'\r\n            - '/TN'\r\n    filter_local:\r\n        CommandLine|contains:\r\n            - '/S localhost'\r\n            - '/S 127.0.0.1'\r\n            - '/S ::1'\r\n    condition: selection_image and selection_create and not filter_local\r\n---\r\ntitle: Remote Scheduled Task Deletion\r\nid: d09ae614-98e5-4cb1-b05a-aab8070bdced\r\nname: remote_schtask_deletion\r\nlogsource:\r\n    product: windows\r\n    category: process_creation\r\ndetection:\r\n    selection_image:\r\n        Image|endswith: '\\schtasks.exe'\r\n    selection_delete:\r\n        CommandLine|contains|all:\r\n            - '/delete'\r\n            - '/S'\r\n            - '/TN'\r\n    filter_local:\r\n        CommandLine|contains:\r\n            - '/S localhost'\r\n            - '/S 127.0.0.1'\r\n            - '/S ::1'\r\n    condition: selection_image and selection_delete and not filter_local\r\n---\r\ntitle: Remote Scheduled Task Deletion After Creation\r\nid: a945aff6-1fc0-4017-90ad-284ffdef18d9\r\nstatus: test\r\ndescription: |\r\n    Detects deletion of scheduled tasks on remote systems shortly after creation, indicating cleanup of lateral movement artifacts.\r\n    This pattern is commonly observed when attackers create tasks for execution and then remove them to cover their tracks.\r\n    Note: This correlation groups only by User, so it may match operations against different target systems by the same account.\r\nreferences:\r\n    - https://attack.mitre.org/techniques/T1070/004/\r\n    - https://attack.mitre.org/techniques/T1053/005/\r\n    - https://redcanary.com/blog/scheduled-tasks/\r\nauthor: SimSpace\r\ndate: 2026-01-19\r\ntags:\r\n    - attack.defense-evasion\r\n    - attack.lateral-movement\r\n    - attack.t1070.004\r\n    - attack.t1053.005\r\ncorrelation:\r\n    type: temporal_ordered\r\n    rules:\r\n        - remote_schtask_creation\r\n        - remote_schtask_deletion\r\n    group-by:\r\n        - User\r\n    timespan: 1h\r\nfields:\r\n    - CommandLine\r\n    - User\r\n    - ParentImage\r\nfalsepositives:\r\n    - Legitimate IT administration using remote task scheduling with immediate cleanup\r\nlevel: high"
query_type: sigma
version: "1"
