# entries/sigma/019c4381-e470-7374-9c6c-15da7d51089c.yaml
id: 019c4381-e470-7374-9c6c-15da7d51089c
name: Scheduled Task Creation with 'TPMProfiler' in Command Line
description: This rule detects the creation of a scheduled task with 'TPMProfiler' in its command line. This specific task name has been observed in malicious contexts related to SolarWinds Web Help Desk exploitation, typically used for persistence to run QEMU for an SSH backdoor. While legitimate TPM-related tasks exist, the 'TPMProfiler' name in this context is indicative of compromise.
query: |
    title: SolarWinds Web Help Desk Exploitation - TPMProfiler Scheduled Task
    id: a1b2c3d4-e5f6-4a5b-8c9d-solarwinds006
    status: experimental
    description: |
        In this attack, threat actors create a scheduled task named TPMProfiler for persistence. While legitimate TPM-related tasks exist, this specific TPMProfiler task name has only been observed in malicious contexts across SolarWinds WHD exploitation. The task typically runs QEMU for an SSH backdoor.
    references:
        - https://www.huntress.com/blog/active-exploitation-solarwinds-web-help-desk-cve-2025-26399
    author: Daniel Jeremiah
    date: 2026/02/07
    tags:
        - attack.persistence
        - attack.t1053.005
        - cve.2025.26399
    logsource:
        category: process_creation
        product: windows
    detection:
        selection:
            CommandLine|contains: 'TPMProfiler'
        condition: selection
    falsepositives:
        - Legitimate TPM-related scheduled tasks
    level: high
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0002
    - TA0004
techniques:
    - T1053.005
tags:
    - attack.persistence
    - attack.t1053.005
    - cve.2025.26399
author: Daniel Jeremiah
version: "2"
references:
    - https://www.huntress.com/blog/active-exploitation-solarwinds-web-help-desk-cve-2025-26399
data_sources:
    - process_creation
    - windows
