# entries/sigma/019c5cfd-4bed-716e-82d2-b7d48cce181e.yaml
id: 019c5cfd-4bed-716e-82d2-b7d48cce181e
name: KongTuke Malware - Scheduled Task SoftwareProtection
description: The rule detects the execution of 'run.exe' where its image path contains 'SoftwareProtectionPlatform' and its command line arguments include both 'SoftwareProtectionPlatform' and 'udp.pyw'. This specific combination is indicative of the KongTuke malware establishing persistence through a scheduled task named 'SoftwareProtection'.
query: |
    title: KongTuke Malware - Scheduled Task SoftwareProtection
    id: a1b2c3d4-e5f6-4a5b-8c9d-0e1f2a3b4c5d
    status: experimental
    description: >-
      The rule detects the execution of 'run.exe' where its image path contains
      'SoftwareProtectionPlatform' and its command line arguments include both
      'SoftwareProtectionPlatform' and 'udp.pyw'. This specific combination is
      indicative of the KongTuke malware establishing persistence through a
      scheduled task named 'SoftwareProtection'.
    references:
      - https://malpedia.caad.fkie.fraunhofer.de/details/js.kongtuke
      - https://isc.sans.edu/diary/32498
      - >-
        https://www.recordedfuture.com/research/tag-124-multi-layered-tds-infrastructure-extensive-user-base
      - >-
        https://trustwave.com/en-us/resources/blogs/spiderlabs-blog/yet-another-nodejs-backdoor-yanb-a-modern-challenge
    author: Daniel Jeremiah
    date: 2026/02/14
    modified: 2026/02/14
    logsource:
      category: process_creation
      product: windows
      service: security
    detection:
      selection_cmd:
        - CommandLine|contains: SoftwareProtectionPlatform
        - CommandLine|contains: udp.pyw
      selection_image:
        - Image|endswith: \run.exe
        - Image|contains: SoftwareProtectionPlatform
      condition: selection_cmd and selection_image
    falsepositives:
      - Legitimate run.exe under path containing SoftwareProtectionPlatform running
        udp.pyw (unlikely)
    level: high
    tags:
      - attack.execution
      - attack.t1059.006
      - attack.persistence
      - attack.t1053.005
      - malware.kongtuke
query_type: sigma
severity: high
tactics:
    - TA0002
    - TA0003
    - TA0004
techniques:
    - T1059.006
    - T1053.005
tags:
    - attack.execution
    - attack.t1059.006
    - attack.persistence
    - attack.t1053.005
    - malware.kongtuke
author: Daniel Jeremiah
version: "5"
references:
    - https://malpedia.caad.fkie.fraunhofer.de/details/js.kongtuke
    - https://isc.sans.edu/diary/32498
    - https://www.recordedfuture.com/research/tag-124-multi-layered-tds-infrastructure-extensive-user-base
    - https://trustwave.com/en-us/resources/blogs/spiderlabs-blog/yet-another-nodejs-backdoor-yanb-a-modern-challenge
data_sources:
    - process_creation
    - windows
