# entries/sigma/019ce888-1b5d-717f-9309-38c9efb6e31b.yaml
id: 019ce888-1b5d-717f-9309-38c9efb6e31b
name: APT33 Tickler Backdoor Persistence via Run Key
description: Detects the persistence mechanism used by the Tickler backdoor, associated with the threat actor APT33. The malware creates a Run registry key with the value name 'SharePoint.exe' to execute its payload on system startup.
query: |
    title: APT33 Tickler Backdoor Persistence via Run Key
    id: 4b10045c-8586-41f7-aa72-fb013b9a09fa
    status: experimental
    description: Detects the persistence mechanism used by the Tickler backdoor,
      associated with the threat actor APT33. The malware creates a Run registry key
      with the value name 'SharePoint.exe' to execute its payload on system startup.
    author: detections.ai
    references:
      - https://www.cloudsek.com/blog/a-threat-actor-landscape-assessment-of-ics-ot-targeting-in-the-2026-iran-us-conflict-and-the-scale-of-the-risk
    date: 2026-03-13
    tags:
      - attack.persistence
      - attack.t1547.001
      - attack.g0064
    logsource:
      category: registry_set # Catches registry value creation/modification events.
      product: windows
    detection:
      selection:
        # Detects the specific value name 'SharePoint.exe' being added to common Run keys.
        TargetObject|contains:
          - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\SharePoint.exe'
          - '\SOFTWARE\Wow6432Node\Microsoft\Windows\CurrentVersion\Run\SharePoint.exe'
      condition: selection
    falsepositives:
      - Unlikely, but a legitimate application could potentially use this name.
    level: high
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0004
techniques:
    - T1547.001
tags:
    - attack.persistence
    - attack.t1547.001
    - attack.g0064
author: detections.ai
version: "3"
references:
    - https://www.cloudsek.com/blog/a-threat-actor-landscape-assessment-of-ics-ot-targeting-in-the-2026-iran-us-conflict-and-the-scale-of-the-risk
data_sources:
    - registry_set
    - windows
