# entries/sigma/019c48f3-eb7f-70d6-9657-ea5a2df56aac.yaml
id: 019c48f3-eb7f-70d6-9657-ea5a2df56aac
name: Scheduled Task Created with Elevated Privilege Parameters
description: This rule detects the creation of scheduled tasks using the `schtasks.exe` utility with specific parameters indicating elevated privileges. Specifically, it looks for tasks created with the `/create` command, running under the `SYSTEM` user account (`/ru SYSTEM`), and with the highest privilege level (`/rl HIGHEST`). This combination of parameters is often leveraged by adversaries for defense evasion, persistence, and privilege escalation by ensuring the scheduled task executes with maximum system privileges.
query: "title: Scheduled Task Created with Elevated Privilege Parameters\r\nid: 7b1350be-03bc-4dd5-886a-53c5a5dafceb\r\nstatus: test\r\ndescription: Detects scheduled tasks created with explicit SYSTEM account and highest privilege level parameters, often used for defense evasion and execution.\r\nreferences:\r\n    - https://attack.mitre.org/techniques/T1053/005/\r\n    - https://attack.mitre.org/techniques/T1548/\r\nauthor: SimSpace\r\ndate: 2026-01-01\r\ntags:\r\n    - attack.privilege-escalation\r\n    - attack.t1053.005\r\n    - attack.t1548\r\nlogsource:\r\n    category: process_creation\r\n    product: windows\r\ndetection:\r\n    selection_schtasks:\r\n        Image|endswith: '\\schtasks.exe'\r\n        CommandLine|contains|all:\r\n            - '/create'\r\n            - '/ru'\r\n            - 'SYSTEM'\r\n            - '/rl'\r\n            - 'HIGHEST'\r\n    condition: selection_schtasks\r\nfalsepositives:\r\n    - Unknown\r\nlevel: high\r\n"
query_type: sigma
severity: high
tactics:
    - TA0004
    - TA0002
    - TA0003
    - TA0005
techniques:
    - T1053.005
    - T1548
    - T1053
tags:
    - attack.privilege-escalation
    - attack.t1053.005
    - attack.t1548
author: SimSpace
version: "1"
references:
    - https://attack.mitre.org/techniques/T1053/005/
    - https://attack.mitre.org/techniques/T1548/
data_sources:
    - process_creation
    - windows
