# entries/sigma/019c9346-7397-71c8-bde7-676b21029a7e.yaml
id: 019c9346-7397-71c8-bde7-676b21029a7e
name: Linux Init.d Script Persistence
description: "Detects creation of malicious init.d scripts for persistence\nTo accompany Hunting rule equivalent at \nhttps://detections.ai/rules/019c91c0-b2e0-759d-a61c-78d08e1d3d99"
query: "title: Linux Init.d Script Persistence\nid: f8e7d6c5-b4a3-2109-8765-4321fedcba98\nstatus: experimental\ndescription: |-\n  Detects creation of malicious init.d scripts for persistence\n  To accompany Hunting rule equivalent at \n  https://detections.ai/rules/019c91c0-b2e0-759d-a61c-78d08e1d3d99\nreferences:\n  - https://attack.mitre.org/techniques/T1037/004/\n  - https://www.elastic.co/security-labs/sequel-on-persistence-mechanisms\nauthor: zendannyy\ndate: 2026-02-21\ntags:\n  - attack.persistence\n  - attack.t1037.004\n  - linux.init_scripts\nlogsource:\n  category: process_creation\n  product: linux\ndetection:\n  selection_script_creation:\n    ParentImage|endswith:\n      - '/bash'\n      - '/sh'\n    Image|endswith:\n      - '/chmod'\n      - '/update-rc.d'\n      - '/systemctl'\n  selection_init_paths:\n    CommandLine|contains:\n      - '/etc/init.d/'\n      - '/etc/systemd/system/'\n  selection_enable_operations:\n    CommandLine|contains:\n      - 'enable'\n      - '+x'\n      - 'defaults'\n  condition: selection_script_creation and selection_init_paths and\n    selection_enable_operations\nfalsepositives:\n  - Software package installations\n  - System service management\nlevel: high\n"
query_type: sigma
severity: medium
tactics:
    - TA0003
techniques:
    - T1037.004
    - T1037
    - T1543.002
tags:
    - attack.persistence
    - attack.t1037.004
    - linux.init_scripts
author: zendannyy
version: "1"
references:
    - https://attack.mitre.org/techniques/T1037/004/
    - https://www.elastic.co/security-labs/sequel-on-persistence-mechanisms
data_sources:
    - process_creation
    - linux
