# entries/sigma/019cf2c3-75bd-74fa-bbc7-89fae9157fe4.yaml
id: 019cf2c3-75bd-74fa-bbc7-89fae9157fe4
name: Slopoly Persistence Scheduled Task
description: Detects the creation of a scheduled task named 'Runtime Broker' that executes files from 'C:\ProgramData\Microsoft\Windows\Runtime\'. This technique is used by the Slopoly malware, deployed by the Hive0163 group, to establish persistence by masquerading as a legitimate Windows service.
query: |
    title: Slopoly Persistence Scheduled Task
    id: 232b738c-e6b4-4157-9199-a6e720d0dd25
    status: experimental
    description: Detects the creation of a scheduled task named 'Runtime Broker'
      that executes files from 'C:\ProgramData\Microsoft\Windows\Runtime\'. This
      technique is used by the Slopoly malware, deployed by the Hive0163 group, to
      establish persistence by masquerading as a legitimate Windows service.
    author: detections.ai
    references:
      - https://www.ibm.com/think/x-force/slopoly-start-ai-enhanced-ransomware-attacks
    date: 2026-03-15
    tags:
      - attack.persistence
      - attack.t1053.005
    logsource:
      category: process_creation
      product: windows
    detection:
      selection:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains|all:
          # Detects the creation of a scheduled task
          - '/create'
          # Specific task name used by Slopoly malware, masquerading as a legitimate service
          - 'Runtime Broker'
          # Specific path where the Slopoly script is deployed
          - 'C:\ProgramData\Microsoft\Windows\Runtime\'
      condition: selection
    falsepositives:
      - Legitimate administrative scripts using the same task name and path,
        although highly unlikely.
    level: high
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0002
    - TA0004
techniques:
    - T1053.005
tags:
    - attack.persistence
    - attack.t1053.005
author: detections.ai
version: "2"
references:
    - https://www.ibm.com/think/x-force/slopoly-start-ai-enhanced-ransomware-attacks
data_sources:
    - process_creation
    - windows
