# entries/sigma/019c48f4-3544-737b-bd26-13da62d38573.yaml
id: 019c48f4-3544-737b-bd26-13da62d38573
name: Script File Dropped in Startup Folder
description: This rule detects the creation of script files (batch, VBS, JS, PowerShell, or shortcut files) within Windows Startup folders. It monitors both user-specific and all-users startup locations. The primary purpose is to identify persistence mechanisms where an attacker or malicious software places a script to execute automatically when the system or user logs on.
query: "title: Script File Dropped in Startup Folder\r\nid: a45c627e-368f-4224-a3e9-29d9f39bc7a4\r\nstatus: experimental\r\ndescription: Detects script files (batch, VBS, JS, PowerShell) created in Windows Startup folders for persistence. Monitors both user-specific and all-users startup locations.\r\nreferences:\r\n    - https://attack.mitre.org/techniques/T1547.001/\r\nauthor: SimSpace\r\ndate: 2026-01-19\r\ntags:\r\n    - attack.persistence\r\n    - attack.t1547.001\r\nlogsource:\r\n    product: windows\r\n    category: file_event\r\ndetection:\r\n    selection_location:\r\n        TargetFilename|contains:\r\n            - '\\Start Menu\\Programs\\Startup\\'\r\n            - '\\Start Menu\\Programs\\StartUp\\'\r\n    selection_extension:\r\n        TargetFilename|endswith:\r\n            - '.bat'\r\n            - '.cmd'\r\n            - '.vbs'\r\n            - '.vbe'\r\n            - '.js'\r\n            - '.jse'\r\n            - '.ps1'\r\n            - '.lnk'\r\n    condition: selection_location and selection_extension\r\nfields:\r\n    - TargetFilename\r\n    - Image\r\n    - User\r\nfalsepositives:\r\n    - Legitimate software installers creating startup scripts\r\n    - IT-deployed automation via Group Policy or configuration management\r\n    - User-created convenience scripts\r\nlevel: medium\r\n"
query_type: sigma
severity: medium
tactics:
    - TA0003
    - TA0004
techniques:
    - T1547.001
tags:
    - attack.persistence
    - attack.t1547.001
author: SimSpace
version: "1"
references:
    - https://attack.mitre.org/techniques/T1547.001/
data_sources:
    - file_event
    - windows
