# entries/sigma/019c48f6-447f-72e3-8bf8-c4478e02e628.yaml
id: 019c48f6-447f-72e3-8bf8-c4478e02e628
name: Lateral Movement via Scheduled Task with Network Logon
description: This rule detects lateral movement via scheduled tasks by correlating three distinct events within a 10-minute window on the same computer and for the same user. First, it identifies a successful network logon (EventID 4624, LogonType 3), excluding machine accounts and anonymous logons. Second, it looks for the creation of a scheduled task (EventID 4698). Third, it detects the execution of suspicious processes (PowerShell, cmd, wscript, cscript) where the parent process is either svchost.exe (with 'Schedule' in its command line) or taskeng.exe/taskhostw.exe/taskhost.exe, indicating execution by the Task Scheduler. This sequence suggests an attacker authenticating remotely, creating a scheduled task, and then immediately executing it.
query: "title: Network Logon\r\nname: network_logon_step_1\r\nid: a1b2c3d4-e5f6-7890-1234-567890abcdef\r\nlogsource:\r\n    product: windows\r\n    service: security\r\ndetection:\r\n    selection:\r\n        EventID: 4624\r\n        LogonType: 3\r\n    filter_machine_accounts:\r\n        TargetUserName|endswith: '$'\r\n    filter_anonymous:\r\n        TargetUserName: 'ANONYMOUS LOGON'\r\n    condition: selection and not 1 of filter_*\r\n---\r\ntitle: Scheduled Task Created\r\nname: scheduled_task_created_step_2\r\nid: f0957e04-5756-45b3-a035-03d92759908f\r\nlogsource:\r\n    product: windows\r\n    service: security\r\ndetection:\r\n    selection:\r\n        EventID: 4698\r\n    condition: selection\r\n---\r\ntitle: Suspicious Process Execution by Task Scheduler\r\nname: scheduled_task_execution_step_3\r\nid: c9ad39bb-a7fe-4fe7-93ce-e87c945edd09\r\nlogsource:\r\n    category: process_creation\r\n    product: windows\r\ndetection:\r\n    selection_image:\r\n        Image|endswith:\r\n            - '\\powershell.exe'\r\n            - '\\pwsh.exe'\r\n            - '\\cmd.exe'\r\n            - '\\wscript.exe'\r\n            - '\\cscript.exe'\r\n    selection_parent_svchost:\r\n        ParentImage|endswith: '\\svchost.exe'\r\n        ParentCommandLine|contains: 'Schedule'\r\n    selection_parent_taskhost:\r\n        ParentImage|endswith:\r\n            - '\\taskeng.exe'\r\n            - '\\taskhostw.exe'\r\n            - '\\taskhost.exe'\r\n    condition: selection_image and (selection_parent_svchost or selection_parent_taskhost)\r\n---\r\ntitle: Lateral Movement via Scheduled Task\r\nid: bf9e662c-b4ee-4d7c-bb8d-1bf3214fee31\r\nstatus: experimental\r\ndescription: |\r\n    Detects lateral movement via scheduled task by correlating network logon, task creation, and immediate execution.\r\n    This three-step kill chain identifies an attacker authenticating remotely, creating a scheduled task, and executing it within minutes.\r\nreferences:\r\n    - https://attack.mitre.org/techniques/T1053/005/\r\n    - https://attack.mitre.org/techniques/T1021/\r\n    - https://www.varonis.com/blog/t1053-scheduled-task\r\nauthor: SimSpace\r\ndate: 2026-01-19\r\ntags:\r\n    - attack.lateral-movement\r\n    - attack.execution\r\n    - attack.persistence\r\n    - attack.t1021\r\n    - attack.t1053.005\r\ncorrelation:\r\n    type: temporal_ordered\r\n    rules:\r\n        - network_logon_step_1\r\n        - scheduled_task_created_step_2\r\n        - scheduled_task_execution_step_3\r\n    group-by:\r\n        - ComputerName\r\n        - TargetUserName\r\n    timespan: 10m\r\nfields:\r\n    - ComputerName\r\n    - TargetUserName\r\n    - IpAddress\r\n    - TaskName\r\n    - Image\r\n    - CommandLine\r\n    - ParentImage\r\nfalsepositives:\r\n    - Legitimate remote administration with immediate task execution\r\n    - Automated deployment systems using network authentication\r\n    - Remote management tools executing tasks immediately after creation\r\nlevel: high"
query_type: sigma
version: "1"
