# entries/sigma/452df9c8-1205-46e0-9ee4-84dd79225779.yaml
id: 452df9c8-1205-46e0-9ee4-84dd79225779
name: Suspicious Registry Queries and Scheduled Task Creation - LummaC2
description: Detects registry queries for software inventory and browser settings as well as suspicious scheduled task creation, techniques observed in LummaC2 malware campaigns. Registry queries target installed software and browser preferences for reconnaissance (T1012). Scheduled tasks ensure persistence (T1053.005)
query: "title: Suspicious Registry Queries and Scheduled Task Creation - LummaC2\nid: 9f8ccdf4-b643-4511-84b4-48cab29725d6\nstatus: experimental\ndescription: \n  Detects registry queries for software inventory and browser settings as well\n  as suspicious scheduled task creation, \n    techniques observed in LummaC2 malware campaigns. Registry queries target installed software and browser preferences \n    for reconnaissance (T1012). Scheduled tasks ensure persistence (T1053.005)\nauthor: Kaung Khant Ko\ndate: 2025-08-25\nreferences:\n  - https://securelist.com/lumma-fake-captcha-attacks-analysis/116274/\n    https://www.certego.net/blog/lummastealer/\ntags:\n  - attack.discovery\n  - attack.t1012\n  - attack.persistence\n  - attack.t1053.005\n  - malware.lumma\nlogsource:\n  product: windows\n  category: registry_event\ndetection:\n  registry_selection:\n    TargetObject|contains:\n      - '\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Uninstall'\n      - '\\Software\\Google\\Chrome\\PreferenceMACs'\n  schtasks_selection:\n    EventID: 1\n    Image|endswith: '\\schtasks.exe'\n    CommandLine|contains:\n      - '/create'\n  condition: registry_selection or schtasks_selection\nfalsepositives:\n  - Legitimate administrative scripts or inventory tools may query these\n    registry keys.\n  - System administrators may legitimately use schtasks.exe for automation.\nlevel: high\n"
query_type: sigma
severity: medium
tactics:
    - TA0007
    - TA0003
tags:
    - attack.discovery
    - attack.t1012
    - attack.persistence
    - attack.t1053.005
    - malware.lumma
author: Kaung Khant Ko
version: "1"
references:
    - https://securelist.com/lumma-fake-captcha-attacks-analysis/116274/ https://www.certego.net/blog/lummastealer/
data_sources:
    - registry_event
    - windows
