# entries/sigma/710131b3-1125-412c-b8c0-5f93748f433e.yaml
id: 710131b3-1125-412c-b8c0-5f93748f433e
name: DarkCloud Stealer - RunOnce Registry Persistence
description: Detects the creation or modification of a registry key under HKEY_CURRENT_USER\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce, a common persistence mechanism used by DarkCloud Stealer.
query: |-
    title: DarkCloud Stealer - RunOnce Registry Persistence
    id: c2d3e4f5-a6b7-c8d9-e0f1-a2b3c4d5e6f7
    status: experimental
    description: Detects the creation or modification of a registry key under HKEY_CURRENT_USER\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce, a common persistence mechanism used by DarkCloud Stealer.
    author: Detections.ai
    date: 2025-06-30
    references:
        - https://unit42.paloaltonetworks.com/darkcloud-stealer-and-obfuscated-autoit-scripting/
    logsource:
        product: windows
        category: registry_event
    detection:
        selection:
            TargetObject|contains: '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\'
        condition: selection
    falsepositives:
        - Legitimate applications or system updates might use RunOnce for temporary execution, but it's less common for persistent installations. Further investigation into the process and the value data is recommended.
    level: high
    tags:
        - attack.persistence
        - attack.t1547.001 # Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder
query_type: sigma
severity: high
tactics:
    - TA0003
    - TA0004
techniques:
    - T1547.001
tags:
    - attack.persistence
    - attack.t1547.001
author: Detections.ai
version: "1"
references:
    - https://unit42.paloaltonetworks.com/darkcloud-stealer-and-obfuscated-autoit-scripting/
data_sources:
    - registry_event
    - windows
