# entries/sigma/75629521-3c5d-456f-b99d-1032e3ef7df0.yaml
id: 75629521-3c5d-456f-b99d-1032e3ef7df0
name: Scheduled task created and deleted fastly (ATexec.py)
description: Detects scenarios where an attacker abuse task scheduler capacities to execute commands or elevate privileges.
query: "title: Scheduled task created and deleted fastly (ATexec.py)\ndescription: Detects scenarios where an attacker abuse task scheduler capacities to execute commands or elevate privileges.\nreferences:\n- https://github.com/mdecrevoisier/EVTX-to-MITRE-Attack/tree/master/TA0002-Execution/T1053.005-Scheduled%20Task\n- https://github.com/SecureAuthCorp/impacket/blob/master/examples/atexec.py\n- https://u0041.co/blog/post/1\n- https://blog.menasec.net/2019/03/threat-hunting-25-scheduled-tasks-for.html\ntags:\n- attack.execution\n- attack.t1053.005 # Scheduled Task/Job: Scheduled Task \nauthor: mdecrevoisier\nstatus: experimental\nlogsource:\n  product: windows\n  service: security\ndetection:\n  selection_create:\n    EventID: 4698\n  selection_delete:\n    EventID: 4699\n  #filter:\n  #  SubjectUserSid: 'S-1-5-18'\n  condition: selection_create > selection_delete | group(Computer, TaskName)\n  timeframe: 5m\nfalsepositives:\n- Rare application activity\nlevel: high"
query_type: sigma
severity: high
tactics:
    - TA0002
    - TA0003
    - TA0004
techniques:
    - T1053.005
tags:
    - attack.execution
    - attack.t1053.005
author: mdecrevoisier
version: "1"
references:
    - https://github.com/mdecrevoisier/EVTX-to-MITRE-Attack/tree/master/TA0002-Execution/T1053.005-Scheduled%20Task
    - https://github.com/SecureAuthCorp/impacket/blob/master/examples/atexec.py
    - https://u0041.co/blog/post/1
    - https://blog.menasec.net/2019/03/threat-hunting-25-scheduled-tasks-for.html
data_sources:
    - windows
    - security
